# etc/supervisord.d/confd_update_slurm.ini
[program:confd_update_slurm]
command=/root/bin/confd_update_slurm.sh
redirect_stderr=true
#stdout_syslog=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
