# etc/supervisord.d/diamond.ini
[program:diamond]
command=/bin/diamond -f 
redirect_stderr=true
stdout_logfile=/var/log/supervisord_diamond.log
stdout_capture_maxbytes=1MB
