# etc/supervisord.d/logstash-forwarder.ini
[program:logstash-forwarder]
command=/opt/logstash-forwarder/bin/logstash-forwarder -config /etc/lumberjack.conf
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log



