# etc/supervisord.d/lumerjack.ini
[program:lumberjack]
command=/opt/logstash-forwarder/bin/logstash-forwarder -config /etc/lumberjack.conf
redirect_stderr=true
stdout_logfile=/var/log/supervisord_lumberjack.log


