# etc/supervisord.d/munge.ini
[program:munge]
command=/sbin/munged -F
redirect_stderr=true
stdout_syslog=true
