# etc/supervisord.d/slurmd.ini
[program:slurmd]
command=/usr/local/sbin/slurmd -D -v
redirect_stderr=true
stdout_syslog=true
