# etc/supervisord.d/slurmd.ini
[program:munge]
command=/sbin/munged -F
redirect_stderr=true
stdout_syslog=true

[program:slurmd]
command=/usr/local/sbin/slurmd -D -v
redirect_stderr=true
stdout_syslog=true


