# etc/supervisord.d/slurmd.ini
[program:slurmd]
command=/root/bin/start_slurmd.sh
redirect_stderr=true
stdout_syslog=true
